# .config/systemd/user/ansible-sync-pull.timer
[Unit]
Description=Periodically pull GitHub ansible-sync repo and converge workstation

[Timer]
OnBootSec=5m
OnUnitActiveSec=90m
RandomizedDelaySec=300
Persistent=true
Unit=ansible-sync-pull.service

[Install]
WantedBy=timers.target
